# deploy/systemd/magic-master.service
[Unit]
Description=Magic Master — автоматический аудио-мастеринг
Documentation=https://github.com/denisok-ai/audio-mastering-web
After=network.target
Wants=network.target
OnFailure=magic-master-alert@%n.service

[Service]
Type=simple
User=magicmaster
Group=magicmaster
WorkingDirectory=/opt/magic-master/backend
Environment="PATH=/opt/magic-master/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/magic-master/.env
ExecStart=/opt/magic-master/venv/bin/python run_production.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Лимиты ресурсов (под сервер 12 GB RAM, 6 vCPU)
LimitNOFILE=65536
LimitNPROC=4096
# 9 GB soft limit (swap pressure), 10 GB hard kill
MemoryHigh=9G
MemoryMax=10G
CPUAffinity=0-5

# Безопасность (ослаблено для совместимости с VPS/контейнерами — 226/NAMESPACE)
# При необходимости можно включить: NoNewPrivileges=true, PrivateTmp=true, ProtectSystem=strict
ReadWritePaths=/opt/magic-master /tmp
ProtectSystem=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
